# archlinux/roles/systemd/files/systemd-journal-gatewayd.socket.d/socket.conf
[Socket]
ListenStream=
ListenStream=/run/journald.sock
