# archlinux/roles/docker/files/docker-prune.timer
[Install]
WantedBy=multi-user.target

[Timer]
OnBootSec=0
OnCalendar=*:0/15
