# ubuntu16.04/roles/tsugi-project/files/timers/tsugi-errors-check.timer
[Unit]
Description=Tsugi Check Errors Timer

[Timer]
OnCalendar=*:0/30
