# archlinux/roles/consul-backup/files/consul-backup.service
[Unit]
Description=Consul Backup

[Service]
ExecStart=/usr/bin/bash /usr/local/bin/consul-backup.sh
