# archlinux/roles/consul-backup/files/consul-backup.timer
[Unit]
Description=Consul Backup Timer

[Timer]
OnBootSec=60
OnCalendar=*:0/15
