# ubuntu16.04/roles/tsugi-project/files/services/tsugi-update-tools.service
[Unit]
Description=Update Tsugi Tools

[Service]
User=www-data
Group=www-data

ExecStart=/bin/bash /usr/local/bin/tsugi-update-tools.sh
