# ubuntu16.04/roles/tsugi-project/files/services/tsugi-update-tools.service
[Unit]
Description=Update Tsugi Tools

[Service]
ExecStart=/bin/bash /usr/local/bin/tsugi_update_tools.sh
