# ubuntu16.04/roles/tsugi-project/files/services/tsugi-errors-check.service
[Unit]
Description=Check Tsugi Logs For Errors

[Service]
ExecStart=/bin/bash /usr/local/bin/tsugi-errors-check.sh
