# ubuntu16.04/roles/tsugi-project/files/services/tsugi-upgrade-database.service
[Unit]
Description=Upgrade Tsugi Database

[Service]
User=www-data
Group=www-data

ExecStart=/bin/bash /usr/local/bin/tsugi-upgrade-database.sh
