# ubuntu16.04/roles/tsugi-project/files/services/tsugi-update.service
[Unit]
Description=Update Tsugi

[Service]
ExecStart=/bin/bash /usr/local/bin/tsugi_update.sh
