# ubuntu16.04/roles/tsugi-project/files/services/tsugi-upgrade-database.service
[Unit]
Description=Upgrade Tsugi Database

[Service]
ExecStart=/bin/bash /usr/local/bin/tsugi_upgrade_database.sh
