# ubuntu16.04/roles/tsugi-project/files/services/tsugi-update-main-page.service
[Unit]
Description=Update Tsugi Main Page

[Service]
ExecStart=/bin/bash /usr/local/bin/tsugi_update_main_page.sh
