# archlinux/roles/unbound/files/unbound.service.d/00-partof-systemd-resolved.conf
[Unit]
Requires=systemd-resolved.service
Before=systemd-resolved.service
PartOf=systemd-resolved.service
After=
